# elasticbeanstalk/supervisord/example.conf
[program:example]
command=
directory=/var/app/current
stdout_logfile=/var/log/example.log
redirect_stderr=true
autostart=true
autorestart=true
